# bashrc/profile.d/umask.sh
# By default, the umask should be set.
umask 022
